# .docker/etc/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command = /usr/sbin/php5-fpm -F
priority = 5
stdout_logfile	= /var/log/supervisor/%(program_name)s.log
stderr_logfile	= /var/log/supervisor/%(program_name)s.log
autorestart	= true